# packages/pkg-anonymize/pkg-anonymize-0.0.1.tar.gz/pkg-anonymize-0.0.1/service/anonymize.service
[Unit]
Description=Anonymize your personnals informations
Wants=network-online.target
Before=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/anonymize execute

[Install]
WantedBy=multi-user.target
